# distrobox/services/apps.service
[Unit]
Description=Preserves flatpak paths to volume

[Service]
ExecStart=/apps.sh

[Install]
WantedBy=default.target
