# distrobox/services/desktop.service
[Unit]
Description=Preserves flatpak paths to volume

[Service]
ExecStart=/desktop.sh

[Install]
WantedBy=default.target
