# build/nvidia-egl-desktop/systemd/entrypoint.service
[Unit]
Description=Systemd port of nvidia-egl-docker's 'entrypoint' supervisor service
Requires=user@1000.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/entrypoint-service.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
